# extra/openrc/cs16-server-1
#!/sbin/openrc-run

description="Counter-Strike 1.6 Dedicated Server"
command="/usr/bin/cs16-server"

start() {
	ebegin "Starting ${description}"
	${command} start > /dev/null 2>&1
}

stop() {
	ebegin "Stopping ${description}"
	${command} stop > /dev/null 2>&1
}
