# extra/openrc/cs16-server-2
#!/sbin/openrc-run

description="Counter-Strike 1.6 Dedicated Server"
command_user="hlds:hlds"
directory="/home/hlds/.steam/cs16"
command="${directory}/hlds_run"

start_pre() {
	while IFS= read -r line || [ -n "${line}" ]; do
		eval "${line}"
	done < /etc/hlds/cs16-server.conf
	command_args="${server_game} ${server_params} ${server_dll}"
}

start() {
	ebegin "Starting ${description}"
	start-stop-daemon ${command} --chdir ${directory} -b -u ${command_user} -- ${command_args}
}

stop() {
	ebegin "Stopping ${description}"
	killall -s SIGINT hlds_run hlds_linux
}
