# cs16-server/cs16-server.service
[Unit]
Description=Counter-Strike 1.6 Dedicated Server
After=network.target

[Service]
Type=simple
User=hlds
KillSignal=SIGINT
EnvironmentFile=/etc/hlds/cs16-server.conf
WorkingDirectory=/var/hlds/cs16
ExecStart=/var/hlds/cs16/hlds_run -game cstrike ${server_params}
StandardOutput=null
StandardInput=null
Restart=on-failure

[Install]
WantedBy=multi-user.target
